# packages/scheduler/gg.mk.innie.opencode-server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>gg.mk.innie.opencode-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{OPENCODE_PATH}}</string>
        <string>serve</string>
        <string>--port</string>
        <string>{{OPENCODE_PORT}}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENCODE_SERVER_PASSWORD</key>
        <string>{{OPENCODE_SERVER_PASSWORD}}</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>NODE_EXTRA_CA_CERTS</key>
        <string>{{NODE_EXTRA_CA_CERTS}}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{MEMORY_DIR}}/logs/opencode-server.log</string>

    <key>StandardErrorPath</key>
    <string>{{MEMORY_DIR}}/logs/opencode-server.error.log</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_DIR}}</string>
</dict>
</plist>
